# deployment/supervisor-workers.conf
; ============================================================
; SecureLicence — Supervisor Configuration
; Place in: /etc/supervisor/conf.d/securelicence.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update
; ============================================================

; ── Queue Worker ────────────────────────────────────────────
; Processes jobs from the database queue (emails, notifications, etc.)
[program:securelicence-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/securelicence/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/securelicence/storage/logs/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

; ── Task Scheduler ──────────────────────────────────────────
; Runs Laravel scheduler every minute (replaces cron for artisan schedule:run)
[program:securelicence-scheduler]
process_name=%(program_name)s
command=php /var/www/securelicence/artisan schedule:work
autostart=true
autorestart=true
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/securelicence/storage/logs/scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
